# files/usr/lib/systemd/user/topgrade.service
[Unit]
Description=Update System
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=/usr/bin/topgrade
