# files/usr/lib/systemd/user/topgrade.timer
[Unit]
Description=Run topgrade 1h after boot

[Timer]
OnBootSec=1h

[Install]
WantedBy=timers.target
